# systemd/sentinel-firmware-update-dispatch.path.in
[Unit]
Description=Watch for Sentinel Firmware Update Requests

[Path]
PathExistsGlob=@STATE_DIR@/requests/*.req
Unit=sentinel-firmware-update-dispatch.service

[Install]
WantedBy=multi-user.target
